# build/base/supervisord.conf
; vibespace supervisor configuration
; Manages: sshd (primary) + ttyd (browser fallback)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; sshd - SSH server for CLI terminal access (primary)
; Port 22 - fast, native terminal protocol
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

; ttyd - Web terminal server (browser fallback)
; Port 7681 - runs as non-root user
; Uses login shell (-l) to source /etc/profile.d/vibespace.sh for env vars
[program:ttyd]
command=/usr/bin/ttyd -W -p 7681 /bin/bash -l
user=user
directory=/vibespace
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd-error.log
environment=HOME="%(ENV_USER_HOME)s",USER="user",PATH="/home/user/.local/bin:/home/user/.npm-global/bin:/usr/local/bin:/usr/bin:/bin",VIBESPACE_ID="%(ENV_VIBESPACE_ID)s",VIBESPACE_NAME="%(ENV_VIBESPACE_NAME)s",VIBESPACE_AGENT="%(ENV_VIBESPACE_AGENT)s",VIBESPACE_SHARE_CREDENTIALS="%(ENV_VIBESPACE_SHARE_CREDENTIALS)s"
