# api/pkg/template/images/supervisord.conf
[supervisord]
nodaemon=true
user=workspace
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=workspace

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Code Server (VS Code in browser)
# Runs on port defined by VIBESPACE_CODE_PORT (default: 8080)
[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:%(ENV_VIBESPACE_CODE_PORT)s --auth none /workspace
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/tmp/code-server.log
stderr_logfile=/tmp/code-server.log
environment=HOME="/home/workspace"

# Preview Server (dev server - template-specific)
# Runs on port defined by VIBESPACE_PREVIEW_PORT (default: 3000)
[program:preview-server]
command=/home/workspace/scripts/preview.sh
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/tmp/preview-server.log
stderr_logfile=/tmp/preview-server.log
environment=HOME="/home/workspace",PORT="%(ENV_VIBESPACE_PREVIEW_PORT)s"

# Production Server (production build server - template-specific)
# Runs on port defined by VIBESPACE_PROD_PORT (default: 3001)
[program:prod-server]
command=/home/workspace/scripts/prod.sh
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/tmp/prod-server.log
stderr_logfile=/tmp/prod-server.log
environment=HOME="/home/workspace",PORT="%(ENV_VIBESPACE_PROD_PORT)s"
