# api/pkg/template/images/supervisord.conf
[supervisord]
nodaemon=true
user=workspace
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=workspace

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Caddy Reverse Proxy (internal routing)
# Listens on port 8080 (single Knative port) and routes to backend services based on Host header
# Priority: 1 (starts first, other services depend on it)
[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile
directory=/workspace
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/workspace"

# Code Server (VS Code in browser)
# Runs on port 8081 (internal, proxied by Caddy from 8080)
# Priority: 10 (starts after Caddy)
[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8081 --auth none /workspace
directory=/workspace
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/workspace"

# Preview Server (dev server - template-specific)
# Runs on port 3000 (internal, proxied by Caddy)
# Priority: 20 (starts after code-server)
[program:preview-server]
command=/home/workspace/scripts/preview.sh
directory=/workspace
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/workspace",PORT="3000"

# Production Server (production build server - template-specific)
# Runs on port 3001 (internal, proxied by Caddy)
# Priority: 30 (starts after preview server)
[program:prod-server]
command=/home/workspace/scripts/prod.sh
directory=/workspace
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/workspace",PORT="3001"
