# api/pkg/image/supervisord.conf
; vibespace supervisor configuration
; Manages: ttyd, port-detector, nats-bridge

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ttyd - Web terminal server
; Runs as non-root user, serves on port 7681
[program:ttyd]
command=/usr/bin/ttyd -W -t fontSize=14 -t fontFamily="JetBrains Mono, monospace" /bin/bash
user=user
directory=/vibespace
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd-error.log
environment=HOME="/home/user",USER="user",PATH="/home/user/.local/bin:/home/user/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"

; port-detector - Monitors listening ports and publishes to NATS
[program:port-detector]
command=/usr/bin/port-detector
user=user
directory=/vibespace
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/port-detector.log
stderr_logfile=/var/log/supervisor/port-detector-error.log
startsecs=5
startretries=10

; nats-bridge - Routes messages between Claude CLI and NATS
[program:nats-bridge]
command=/usr/bin/nats-bridge
user=user
directory=/vibespace
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nats-bridge.log
stderr_logfile=/var/log/supervisor/nats-bridge-error.log
startsecs=5
startretries=10
